# deploy/linux/rosina-remote.service
[Unit]
Description = rosina service

[Service]
Type = simple
Environment = "VSOCKET_ADDR=:12001"
Environment = "VPIN_1=18"
Environment = "VPIN_2=16"
Environment = "VSWITCH_DURATION = 300"
ExecStart = /usr/local/rosina
Restart = on-abort
StartLimitInterval = 60
StartLimitBurst = 10

[Install]
WantedBy = multi-user.target

